# mid-server/asset/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:mid]
command=/opt/agent/bin/mid.sh start
user=root
stopsignal=TERM
stopwaitsecs=30
stopasgroup=true